# pg-backup/crontab
0 2 * * * /backup.sh >> /var/log/pg-backup.log 2>&1
